# files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI

# Wireless radio0 specific settings
set wireless.radio0.htmode='HE40'
set wireless.radio0.channel='auto'
set wireless.radio0.country='US'
set wireless.radio0.itxbfen='1'
set wireless.radio0.legacy_rates='0'
set wireless.default_radio0.ocv='0'
set wireless.default_radio0.rts='512'
set wireless.default_radio0.short_preamble='1'
set wireless.default_radio0.dtim_period='3'
set wireless.default_radio0.short_gi_20='1'
set wireless.default_radio0.short_gi_40='1'
set wireless.default_radio0.short_gi='1'
set wireless.default_radio0.uapsd='0'

# Wireless radio1 specific settings
set wireless.radio1.htmode='HE160'
set wireless.radio1.channel='auto'
set wireless.radio1.country='US'
set wireless.radio1.itxbfen='1'
set wireless.default_radio1.ocv='0'

# Firewall settings
set firewall.@zone[1].input='DROP'
set firewall.@zone[1].forward='DROP'

# uHTTPd settings
set uhttpd.main.redirect_https='1'

EOI
